# roles/btrbk/files/etc/systemd/system/btrbk.timer
[Unit]
Description=btrbk backup

[Timer]
OnCalendar=*:0/30
Persistent=true
AccuracySec=1min

[Install]
WantedBy=timers.target
