# roles/miscellaneous/files/etc/systemd/system/btrfsmaintenance-refresh.path
[Unit]
Description=Watch /etc/default/btrfsmaintenance

[Path]
PathChanged=/etc/default/btrfsmaintenance

[Install]
WantedBy=multi-user.target
